# configs/autostart/pipewire.desktop
[Desktop Entry]
Type=Application
Name=PipeWire
Comment=PipeWire Multimedia Service
Exec=pipewire
Hidden=false
NoDisplay=true
X-KDE-autostart-phase=1
X-KDE-autostart-after=panel
