# configs/autostart/steam.desktop
[Desktop Entry]
Name=Steam
Exec=steam -silent
Type=Application
NoDisplay=true
