# configs/autostart/wireplumber.desktop
[Desktop Entry]
Type=Application
Name=WirePlumber
Comment=PipeWire Session Manager
Exec=/bin/sh -c 'sleep 1 && wireplumber'
Hidden=false
NoDisplay=true
X-KDE-autostart-phase=1
